# private/supervisor/gunicorn.conf
[program:gunicorn]
command=/root/venv/bin/gunicorn -w 2 -b :5000 fisher:app
directory=/home/ubuntu/fisher
autostart=true
autorestart=true
stopasgroup=true
user=root
redirect_stderr=true
